# recipes-support/libusbgx/files/usbg.service
[Unit]
Description=Load default USB gadget schema g1.schema
ConditionPathExists=/etc/usbg/g1.schema
Before=systemd-networkd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/gadget-import g1 /etc/usbg/g1.schema
ExecStartPost=/bin/sh -c '/bin/echo $((0x4000+$(expr $(cat /proc/device-tree/toradex,product-id) + 0))) > \
                                    /sys/kernel/config/usb_gadget/g1/idProduct'
ExecStartPost=/bin/sh -c '/bin/cat /proc/device-tree/serial-number > \
                                   /sys/kernel/config/usb_gadget/g1/strings/0x409/serialnumber'
ExecStartPost=/bin/sh -c '/bin/echo `ls /sys/class/udc/` > \
                                    /sys/kernel/config/usb_gadget/g1/UDC'

[Install]
WantedBy=multi-user.target
