# recipes-core/tezi-init/files/rc.local
#!/bin/sh

# Start services and customize the boot process here.
echo "Running /etc/rc.local..."

# Start Toradex Installer Qt user interface in background
/usr/bin/tezi -qws &

